# curator1-Margaritka.conf
[program:curator1-Margaritka]
command=python3 /home/curator1/curator1-Margaritka/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator1/curator1-Margaritka/test.err.log
stdout_logfile=/home/curator1/curator1-Margaritka/test.out.log
